# deployment/agent.service
[Unit]
Description=NekoProxy Agent Service
After=network.target wg-quick@wg0.service
Wants=wg-quick@wg0.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/nekoproxy
Environment="PATH=/opt/nekoproxy/venv/bin"
Environment="NEKO_AGENT_CONTROLLER_URL=http://10.0.0.1:8000"
Environment="NEKO_AGENT_WIREGUARD_IP=10.0.0.2"
ExecStart=/opt/nekoproxy/venv/bin/python -m agent.main
Restart=always
RestartSec=5

# Agent needs root for binding to privileged ports
# Or use CAP_NET_BIND_SERVICE capability
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
