# deployment/controller.service
[Unit]
Description=NekoProxy Controller Service
After=network.target

[Service]
Type=simple
User=nekoproxy
Group=nekoproxy
WorkingDirectory=/opt/nekoproxy
Environment="PATH=/opt/nekoproxy/venv/bin"
Environment="NEKO_DATABASE_URL=sqlite:////var/lib/nekoproxy/controller.db"
ExecStart=/opt/nekoproxy/venv/bin/python -m controller.main
Restart=always
RestartSec=5

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/nekoproxy /var/log/nekoproxy
PrivateTmp=true

[Install]
WantedBy=multi-user.target
